# services/cups-browsed/run
#!/usr/bin/execlineb -P
/usr/sbin/cups-browsed
