# services/cups/run
#!/usr/bin/execlineb -P
/usr/sbin/cupsd -f
