# services/avahi-daemon/run
#!/usr/bin/execlineb -P
/usr/sbin/avahi-daemon
